# roles/letsencrypt/files/cron.daily
#!/bin/sh

/usr/sbin/dehydrated -c -k /etc/letsencrypt/hooks.sh

